# vm/cowsay.service
[Unit]
Description=
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cowsay-service

[Install]
WantedBy=multi-user.target
